# systemd-web.service
[Unit]
Description=Embedded Systemd Web Manager
After=network.target

[Service]
Type=simple
# Run as root so it has permissions to execute systemctl start/stop/restart
User=root
# Execute the binary with our new flags
ExecStart=/usr/local/bin/systemd-web --bind=127.0.0.1 --port=6999
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target